# root/app/crontabs/daily
0 0 * * * /app/gen/daily >> /config/log/snapshot/daily.log