# root/app/crontabs/hourly
0 * * * * /app/gen/hourly >> /config/log/snapshot/hourly.log