# root/app/crontabs/weekly
0 0 * * 0 /app/gen/weekly >> /config/log/snapshot/weekly.log