# root/etc/services.d/crond/run
#!/usr/bin/with-contenv bash
crond -f -L /config/log/cron.log